# lib/crontab
# Edit this file to introduce tasks to be run by cron.
# #
SHELL=/bin/bash
33 9 * * * bash -l -c '$HOME/.bin/djsd -d'
